# infrastructure/systemd/api-gateway.service
[Unit]
Description=Crypto Trading Platform - API Gateway
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Workspace/crypto-trading-platform/apps/api-gateway
Environment="NODE_ENV=production"
EnvironmentFile=/home/ubuntu/Workspace/crypto-trading-platform/apps/api-gateway/.env
ExecStart=/usr/bin/node /home/ubuntu/Workspace/crypto-trading-platform/apps/api-gateway/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=api-gateway
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
