# infrastructure/systemd/bot-orchestrator.service
[Unit]
Description=Crypto Trading Platform - Bot Orchestrator
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Workspace/crypto-trading-platform/apps/bot-orchestrator
Environment="NODE_ENV=production"
EnvironmentFile=/home/ubuntu/Workspace/crypto-trading-platform/apps/bot-orchestrator/.env
ExecStart=/usr/bin/node /home/ubuntu/Workspace/crypto-trading-platform/apps/bot-orchestrator/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot-orchestrator
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
