# data/bot-instances/Js1Gaz4sMPPiDNgFbmAgDFLe4je2/Js1Gaz4sMPPiDNgFbmAgDFLe4je2-pool-1/supervisor/bot-anshjarvis2003-bot-1.conf
[program:bot-anshjarvis2003-bot-1]
command=freqtrade trade --config /pool/bots/anshjarvis2003-bot-1/config.json --strategy-path /pool/strategies --strategy EmaRsiStrategy --db-url sqlite:////pool/bots/anshjarvis2003-bot-1/tradesv3.sqlite --logfile /pool/bots/anshjarvis2003-bot-1/freqtrade.log
directory=/freqtrade
user=ftuser
autostart=true
autorestart=true
startretries=3
startsecs=10
redirect_stderr=true
stdout_logfile=/pool/logs/bot-anshjarvis2003-bot-1.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=FREQTRADE_USER_DATA_DIR="/pool/bots/anshjarvis2003-bot-1"
