# apps/bot-orchestrator/bot-manager.service
[Unit]
Description=Freqtrade Bot Manager
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Workspace/crypto-trading-platform/bot-manager
EnvironmentFile=/home/ubuntu/Workspace/crypto-trading-platform/bot-manager/.env
ExecStart=/usr/bin/env node /home/ubuntu/Workspace/crypto-trading-platform/bot-manager/index.js

# Always restart on failure or exit
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bot-manager

# Health check
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
