# apps/bot-orchestrator/portfolio-monitor.service
[Unit]
Description=Portfolio Monitor Service - Continuous FreqTrade Portfolio Monitoring
Documentation=https://github.com/Git-Ansh/Crypto-Pilot-Freqtrade
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
Group=root
WorkingDirectory=/home/ubuntu/Workspace/crypto-trading-platform/bot-manager
ExecStart=/usr/bin/node portfolio-monitor-service.js
ExecReload=/bin/kill -HUP $MAINPID

# Environment variables
Environment=NODE_ENV=production
Environment=BOT_BASE_DIR=/home/ubuntu/Workspace/crypto-trading-platform/freqtrade-instances

# Restart configuration
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# Resource limits
LimitNOFILE=8192
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-monitor

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/Workspace/crypto-trading-platform/freqtrade-instances

# Health checks
WatchdogSec=300

[Install]
WantedBy=multi-user.target
Alias=portfolio-monitor.service
