# apps/bot-orchestrator/sync-service.service
[Unit]
Description=FreqTrade Local-to-Turso Database Sync Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Workspace/crypto-trading-platform/bot-manager
Environment=NODE_ENV=production
Environment=PATH=/home/ubuntu/.turso:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node sync-service.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment file
EnvironmentFile=-/home/ubuntu/Workspace/crypto-trading-platform/bot-manager/.env

[Install]
WantedBy=multi-user.target
