# apps/bot-orchestrator/turso-token-scheduler.service
[Unit]
Description=Turso Token Auto-Renewal Scheduler
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Workspace/crypto-trading-platform/bot-manager
ExecStart=/usr/bin/node refresh-turso-token.js --schedule
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
